# common/hot/softwareconfig/rc.local
#!/bin/bash
nohup /usr/bin/heat-agent &>/dev/null &
